# ubuntu/couchbase/couchbase-wsoc/build/supervisord.conf
 
[supervisord]
nodaemon=true

[program:couchbase]
command=su - couchbase -c "/opt/couchbase/bin/couchbase-server -- -noinput"

[program:websocketsApp]
command=/usr/local/bin/main